# docker/supervisor/conf.d/laravel-scheduler.conf
[program:laravel-scheduler]
process_name=%(program_name)s
command=php /var/www/html/artisan schedule:work
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
startsecs=0
stopwaitsecs=10
stopsignal=TERM
numprocs=1
priority=999
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
